# deploy/systemd/user/sage-mate-tunnel.service
[Unit]
Description=Sage Mate Cloudflare tunnel
After=sage-mate-site.service
Requires=sage-mate-site.service

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__
Environment=TUNNEL_PROTOCOL=http2
ExecStart=__REPO_ROOT__/tools/run_named_tunnel.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
