# deploy/systemd/user/sage-faculty-twin-vllm-openai-proxy.service
[Unit]
Description=Sage Faculty Twin OpenAI-compatible vLLM proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__
Environment=PYTHON_BIN=__PYTHON_BIN__
Environment=VLLM_PROXY_HOST=127.0.0.1
Environment=VLLM_PROXY_PORT=18001
Environment=VLLM_PROXY_UPSTREAM_BASE_URL=http://127.0.0.1:18000/v1
Environment=VLLM_PROXY_PATH_PREFIX=/v1
ExecStart=__REPO_ROOT__/tools/run_vllm_openai_proxy.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
