# deploy/systemd/user/sage-mate-inference-monitor.service
[Unit]
Description=Check and repair Sage Mate inference path
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__REPO_ROOT__
Environment=PYTHON_BIN=__PYTHON_BIN__
ExecStart=/bin/bash __REPO_ROOT__/tools/monitor_twin_inference.sh
