# deploy/systemd/user/sage-mate-site.service
[Unit]
Description=Sage Mate local site proxy
After=sage-mate-app.service
Wants=sage-mate-app.service

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__
Environment=APP_PORT=55601
Environment=HOMEPAGE_REDIRECT_ORIGIN=https://me.sage.org.ai
Environment=SITE_PORT=8088
ExecStart=__REPO_ROOT__/tools/run_local_proxy.sh
ExecReload=__REPO_ROOT__/tools/reload_local_proxy.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
