# deploy/systemd/user/sage-mate-wiki-sync.timer
[Unit]
Description=Timer for wiki KB sync (every 30 minutes)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
